# com.telegram.bot.plist
# ุญู ูุดููุฉ ุชููู ุงูุจูุช (Bot 24/7 Solution)

## ๐จ ุงููุดููุฉ ุงูุญุงููุฉ

ุงูุจูุช ูุชููู ุจุณุจุจ ุฎุทุฃ **Conflict**:
```
telegram.error.Conflict: Conflict: terminated by other getUpdates request
```

ูุฐุง ูุนูู ุฃู ููุงู **ูุณุฎุฉ ุฃุฎุฑู ูู ุงูุจูุช ุชุนูู** ูู ููุงู ุขุฎุฑ ูุชุญุงูู ุงูุงุชุตุงู ุจููุณ Bot Token.

## ๐ ุฃุณุจุงุจ ูุญุชููุฉ:

1. **Bot ููุดุชุบู ุนูู ุณูุฑูุฑ ุขุฎุฑ** (ูุซู Render, Railway, Heroku, VPS)
2. **Webhook ูุถุจูุท** ุนูู server ุขุฎุฑ (ุชุญูู ูู ุญุณุงุจ Telegram Bot)
3. **ูุณุฎ ูุชุนุฏุฏุฉ ูุญููุงู** (ุชู ุญููุง ุจุงูู scripts ุงูุฌุฏูุฏุฉ)

## โ ุงูุญููู ุงููุทุจูุฉ

ุชู ุฅูุดุงุก ุงููููุงุช ุงูุชุงููุฉ ูุถูุงู ุนูู ุงูุจูุช 24/7:

### 1. Scripts ุงูุชุดุบูู

| Script | ุงููุตู |
|--------|-------|
| `run_bot.sh` | ุชุดุบูู ุงูุจูุช + Webhook + Watchdog |
| `stop_bot.sh` | ุฅููุงู ุฌููุน ุงูุนูููุงุช |
| `status.sh` | ูุญุต ุญุงูุฉ ุงูุจูุช |
| `watchdog.sh` | ูุฑุงูุจุฉ ูุฅุนุงุฏุฉ ุชุดุบูู ุชููุงุฆูุฉ |
| `delete_webhook.py` | ุญุฐู Webhook ูู Telegram |

### 2. ููููุฉ ุงูุงุณุชุฎุฏุงู

```bash
# ุชุดุบูู ุงูุจูุช
./run_bot.sh

# ูุญุต ุงูุญุงูุฉ
./status.sh

# ุฅููุงู ุงูุจูุช
./stop_bot.sh

# ุญุฐู Webhook (ุฅุฐุง ูุงู ููุงู conflict)
python3 delete_webhook.py
```

### 3. ูููุงุช ุงูู Logs

ุฌููุน ุงูุณุฌูุงุช ูุชููุฑุฉ ูู ูุฌูุฏ `logs/`:
- `logs/bot.log` - ุณุฌู ุงูุจูุช ุงูุฑุฆูุณู
- `logs/webhook.log` - ุณุฌู webhook server
- `logs/watchdog.log` - ุณุฌู ุงููุฑุงูุจ ุงูุชููุงุฆู
- `logs/supervisor.log` - ุณุฌู ุงูุชุดุบูู

## ๐ง ุญู ูุดููุฉ Conflict

ุฅุฐุง ุงุณุชูุฑุช ูุดููุฉ Conflict:

### ุงูุฎุทูุฉ 1: ุชุญูู ูู Webhook
```bash
python3 delete_webhook.py
```

### ุงูุฎุทูุฉ 2: ุงุจุญุซ ุนู Deployments ุฃุฎุฑู
ุชุญูู ูู:
- Render (https://render.com)
- Railway (https://railway.app)
- Heroku (https://heroku.com)
- ุฃู VPS ุฃู Server ุขุฎุฑ

### ุงูุฎุทูุฉ 3: ุฃููู ุฌููุน ุงููุณุฎ
```bash
./stop_bot.sh
sleep 5
./run_bot.sh
```

## ๐ฏ ููุนูู 24/7 ุจุฏูู ุชููู

### ุงูุญู ุงูุฃูุซู: ุงุณุชุฎุฏุงู MacOS launchd

<function_calls>
<invoke name="write_to_file">
<parameter name="CodeContent"><?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.telegram.bot</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ayoubelmessaoudi/Downloads/4/Cod Conferme/telegram_digital_bot/run_bot.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/ayoubelmessaoudi/Downloads/4/Cod Conferme/telegram_digital_bot</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/ayoubelmessaoudi/Downloads/4/Cod Conferme/telegram_digital_bot/logs/launchd.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/ayoubelmessaoudi/Downloads/4/Cod Conferme/telegram_digital_bot/logs/launchd_error.log</string>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
